# shairport_classic/rootfs/etc/services.d/shairport/run
#!/usr/bin/with-contenv bash
echo "🔊 Starting Shairport Sync..."
exec shairport-sync -vv -a "HomeAssistant Speaker" -o alsa -- -d hw:1,0
