# shairport_classic/services.d/avahi/run
#!/usr/bin/with-contenv bash
echo "🟢 Starting Avahi & DBus..."

mkdir -p /run/dbus
dbus-daemon --system &
sleep 1
exec avahi-daemon --no-drop-root --debug