# shairport_classic/rootfs/etc/services.d/avahi/run
#!/usr/bin/with-contenv bash
echo "🟢 Starting Avahi..."
exec avahi-daemon --no-chroot --debug
