# shairport-addon/rootfs/etc/services.d/spotifyd/finish
#!/usr/bin/with-contenv bash
echo "Encerrando Shairport Sync"
